# systemd.socket.tpl
[Socket]
ListenStream={{ip}}:{{port}}
NoDelay=true

[Install]
WantedBy=sockets.target

